# root/etc/inittab
## busybox<->docker integration ##
::ctrlaltdel:/bin/kill -QUIT 1
::restart:/bin/false
##################################

::respawn:/usr/bin/gosu syncthing:syncthing /usr/bin/syncthing -gui-address='0.0.0.0:8384'
